# data/xdg-desktop-portal-pantheon.service.in
[Unit]
Description=Portal service (Pantheon implementation)
Requires=io.elementary.settings-daemon.xdg-desktop-portal.service
After=io.elementary.settings-daemon.xdg-desktop-portal.service

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.pantheon
ExecStart=@LIBEXECDIR@/xdg-desktop-portal-pantheon
