# data/xdg-desktop-portal-pantheon.service.in
[Unit]
Description=Portal service (Pantheon implementation)
After=graphical-session.target
Requisite=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.pantheon
ExecStart=@LIBEXECDIR@/xdg-desktop-portal-pantheon
