# data/xdg-desktop-portal-pantheon.service.in
[Unit]
Description=Portal service (Pantheon implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.pantheon
ExecStart=@libexecdir@/xdg-desktop-portal-pantheon
